# fsy-console.service
[Unit]
Description=FSY Packet Console - APRS TNC and Web UI
After=network.target
Wants=network-online.target

[Service]
Type=forking
# User will be current user who installs the service
User=pi
# WorkingDirectory will be the project directory
WorkingDirectory=/home/pi/fsy-packet-console

# Start the console in a screen session with bash wrapper
# This ensures virtual environment is activated
ExecStart=/usr/bin/screen -dmS fsy-console /bin/bash -c 'cd /home/pi/fsy-packet-console && source /home/pi/conole-venv/bin/activate && exec python main.py'

# Stop gracefully (send Ctrl+C to screen session)
ExecStop=/usr/bin/screen -S fsy-console -X quit

# Restart on failure
Restart=on-failure
RestartSec=10

# Cleanup on stop
KillMode=process

# Logging to syslog
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fsy-console

# Process limits (optional but recommended)
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
